# ops/callbacks/coffeeos-callback-audit-alert.service
[Unit]
Description=CoffeeOS callback audit alert check
After=network.target

[Service]
Type=oneshot
User=deploy
Group=deploy
Environment=RAILS_ENV=production
Environment=APP_DIR=/var/www/coffeeos/current
Environment=WINDOW_MINUTES=60
Environment=MAX_ERROR_RATE_PERCENT=2.0
Environment=MIN_EVENTS=20
WorkingDirectory=/var/www/coffeeos/current
ExecStart=/usr/bin/env bash /var/www/coffeeos/current/ops/callbacks/run_callback_audit_alert.sh
